# templates/container.service
[Service]
ExecStart=/usr/bin/systemd-nspawn -bD /var/lib/container/{{ container_name }} --network-bridge={{ bridge_ifname }}
KillMode=process
